# crontab
# m h  dom mon dow   command
11  */2  *  *  *  ${HOME}/.rc/bin/machine-status | ${HOME}/.rc/bin/notify
