# crontab
SHELL=/bin/bash
PATH=${PATH}:${HOME}/bin

# m  h   dom mon dow   command
30  */2  *   *   *     machine-status | slack
